# .docker/build/php_content/logging/supervisor/collect_civicrm_logs.conf
[program:collect_civicrm_logs]
command=tail -n 0 -F /var/www/default/htdocs/httpdocs/sites/default/files/civicrm/ConfigAndLog/CiviCRM.log
numprocs=1
priority=15
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
