# .docker/build/php_content/logging/supervisor/start_rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
numprocs=1
priority=150
autostart=true
autorestart=true
redirect_stderr=true
